# sentinel.timer
[Unit]
Description=Run SentinelOps every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=sentinel.service

[Install]
WantedBy=timers.target
