# sentinel.service
[Unit]
Description=SentinelOps Self-Healing Linux Agent

[Service]
Type=oneshot
ExecStart=/opt/sentinelops/sentinel.sh
